# ext/without-clash-daemon.service
[Unit]
Description=without-clash daemon

[Service]
Type=simple
ExecStart=/usr/bin/without-clash --daemon

[Install]
WantedBy=multi-user.target
